# avahi_proxy.rc
#!/bin/ksh

daemon="/usr/local/bin/avahi-proxy"
daemon_flags="-p 5300"
daemon_user="nobody"

. /etc/rc.d/rc.subr

rc_cmd $1
